# units/ajbridge@.service
[Unit]
Description=Run a2jmidi_bridge on the given port

[Service]
Type=simple
ExecStart=/usr/bin/a2jmidi_bridge "%I"

