# units/makereq.path
[Path]
PathModified=/usr/lib/oslp/json

[Install]
WantedBy=basic.target
