# units/sgwait@.service
[Unit]
Description=Wait for syngrp@%i
PartOf=syngrp@%i.service

[Service]
Type=oneshot
ExecStart=/bin/sleep inf

