# units/xvnc.socket
[Unit]                        
Description=VNC socket

[Socket]
ListenStream=0.0.0.0:5900
Accept=false

[Install]
WantedBy=default.target

